# work15_backup/borg-backup.timer
# /etc/systemd/system/borg-backup.timer
[Unit]
Description=Borg Backup
[Timer]
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
